# cec-controller.service
[Unit]
Description=CEC Controller Service
Documentation=https://github.com/eliottness/cec-controller
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/cec-controller -status
Restart=on-failure
RestartSec=10s

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cec-controller

[Install]
WantedBy=multi-user.target
